# 7.5-domain/systemd/lego-hs74th.timer
[Unit]
Description=lego renew
After=network.target

[Timer]
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
